# tg-bot.service
[Unit]
Description=Telegram AI Planner Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ai-planner-bot
EnvironmentFile=/home/ubuntu/ai-planner-bot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/ai-planner-bot/venv/bin/python -u /home/ubuntu/ai-planner-bot/app.py
Restart=always
RestartSec=5
TimeoutStartSec=60
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tg-bot

[Install]
WantedBy=multi-user.target